# iac/ansible/roles/supervisor/templates/laravel-worker.conf.j2
[program:{{ supervisor_worker_process_name }}]
process_name=%(program_name)s_%(process_num)02d
command={{ supervisor_worker_command }}
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user={{ app_user }}
numprocs={{ supervisor_numprocs }}
redirect_stderr=true
stdout_logfile={{ app_root }}/storage/logs/worker.log
stopwaitsecs=3600
